# socks-vpn.service
[Unit]
Description=SOCKS VPN via tun2socks
After=network-online.target socks-vpn-server.service
Wants=socks-vpn-server.service

[Service]
Type=simple
ExecStart=/usr/local/bin/socks-vpn-control daemon
ExecStop=/usr/local/bin/socks-vpn-control stop
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
