# socks-vpn.service
[Unit]
Description=SOCKS VPN via tun2socks (managed)
After=network-online.target ssh-socks-managed.service
Wants=ssh-socks-managed.service

[Service]
Type=simple
ExecStart=/usr/local/bin/socks-vpn-control daemon
ExecStop=/usr/local/bin/socks-vpn-control stop
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
